# kirkstone/meta-qt5-app/recipes-qt5/system-dashboard/files/dashboard-eglfs.service
[Unit]
Description=Dashboard Application (EGLFS)
After=network.target

[Service]
ExecStart=/usr/bin/dashboard --platform eglfs
Restart=always

[Install]
WantedBy=multi-user.target 