# kirkstone/meta-qt5-app/recipes-qt5/system-dashboard/files/system-dashboard.service
[Unit]
Description=System Dashboard Application
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Environment=QT_QPA_PLATFORM=offscreen
Environment=QT_QPA_FONTDIR=/usr/share/fonts
Environment=QT_QUICK_BACKEND=software
Environment=XDG_RUNTIME_DIR=/tmp/runtime-root
ExecStart=/usr/bin/system-dashboard
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 