# kirkstone/local/update-service/services/update-service.service
[Unit]
Description=Update Service (RAUC Broker)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/update-service
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target