# kirkstone/local/dashboard/services/dashboard.service
[Unit]
Description=Dashboard Application (LinuxFB)
After=network.target

[Service]
ExecStart=/usr/bin/dashboard
Restart=always

[Install]
WantedBy=multi-user.target 